# ansible/templates/system_health_archive.service.j2
[Unit]
Description=Facility Archive Health Monitor
After=network.target

[Service]
Type=simple
User=root
ExecStart=/opt/sysmon/monitor.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
